# target/linux/layerscape/base-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
echo 378 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio378/direction
echo 0 > /sys/class/gpio/gpio378/value
exit 0
